# supervisor.ini
[program:calculation-job-wrapper-correlation]
user=root
redirect_stderr=true
stdout_logfile=syslog
directory=/root/calculation-job-wrapper
command=/root/calculation-job-wrapper/calculation-job-wrapper --config=config.correlation.json
autorestart=true
autostart=true
startretries = 30
startsecs = 10

[program:calculation-job-wrapper-indoorclimate]
user=root
redirect_stderr=true
stdout_logfile=syslog
directory=/root/calculation-job-wrapper
command=/root/calculation-job-wrapper/calculation-job-wrapper --config=config.indoorclimate.json
autorestart=true
autostart=true
startretries = 30
startsecs = 10
